# rootfs/etc/services.d/rtlsdr-airband/run
#!/usr/bin/with-contenv bash
#shellcheck shell=bash

echo "Starting RTLSDR-Airband"
s6-setuidgid root /usr/local/bin/rtl_airband -c /run/rtlsdr-airband/rtl_airband.conf
