# rootfs/etc/services.d/log-monitor/run
#!/usr/bin/with-contenv bash
#shellcheck shell=bash

tail -Fn0 /var/log/daemon.log | \
while read line ; do
        echo "$line"
done