# rootfs/etc/services.d/rsyslsogd/run
#!/usr/bin/with-contenv sh
# shellcheck shell=sh

 service rsyslog start > /dev/null 2>&1
 sleep 600