# rootfs/etc/services.d/icecast/run
s6-setuidgid root icecast -c /etc/icecast.xml